# systemd/stereoio-thecurrent.service
[Unit]
Description=Stereio - The Current
Wants=network-online.target
After=network-online.target

[Service]
Environment=VIRTUAL_ENV=/home/pi/env
Environment=PATH=/home/pi/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/mpv --no-video --input-ipc-server=/tmp/stereoio-thecurrent-socket --mute=yes https://current.stream.publicradio.org/kcmp.mp3

WorkingDirectory=/home/pi/
StandardOutput=inherit
StandardError=inherit
Restart=always
User=jgrover

[Install]
WantedBy=multi-user.target